# deploy/tandem-ws.service
[Unit]
Description=Tandem — y-websocket server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tandem
Environment=NODE_ENV=production
EnvironmentFile=/root/tandem/.env
ExecStart=/usr/local/bin/node server/dist/server/ws.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
