# deploy/tandem-web.service
[Unit]
Description=Tandem — Next.js app
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tandem
Environment=NODE_ENV=production
EnvironmentFile=/root/tandem/.env
ExecStart=/usr/local/bin/node node_modules/next/dist/bin/next start -p 3742
Restart=on-failure
RestartSec=3
# Forward console output to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
